# s6/debian-root/etc/services.d/lighttpd-access-log/run
#!/usr/bin/execlineb -P
cat /var/log/lighttpd/access.log 2>&1