# s6/debian-root/etc/services.d/lighttpd-error-log/run
#!/usr/bin/execlineb -P
cat /var/log/lighttpd/error.log 2>&1
